# crontab
0 12 * * * /usr/bin/certbot renew --post-hook "systemctl reload nginx"
